# nws/nws-notify.rules
SUBSYSTEM=="net",  RUN+="/etc/network-daemon/scripts/nw_notify $env{ACTION} $env{INTERFACE}"
